# build/systemd/maradns.service
# Place this file here:
# /etc/systemd/system/maradns.service
# Then
# systemctl enable maradns
# systemctl start maradns
# To view log
# journalctl -u maradns
[Unit]
 After=network-online.target
 ConditionPathExists=/usr/local/bin/maradns
 Description=MaraDNS
 Wants=network-online.target

[Service]
 ExecStart=/usr/local/bin/maradns
 StandardOutput=journal
 TimeoutSec=0
 Type=exec

[Install]
 WantedBy=multi-user.target
